# deploy/bet2invest-poster.service
[Unit]
Description=Bet2Invest Poster - Service de publication automatique de pronostics
After=network-online.target
Wants=network-online.target

[Service]
# Type=notify : AddSystemd() dans Program.cs implémente le protocole sd_notify
# systemd est informé quand le service est réellement prêt (et non au simple fork)
Type=notify
User=bet2invest
Group=bet2invest
WorkingDirectory=/opt/bet2invest-poster
ExecStart=/usr/local/bin/dotnet /opt/bet2invest-poster/Bet2InvestPoster.dll

# Redémarrage automatique (NFR1 : délai < 30s)
# Restart=always : redémarre aussi après un arrêt propre (exit 0) pour garantir le service continu
Restart=always
RestartSec=5

# Credentials chargées depuis le fichier d'env (NFR6 - jamais dans appsettings.json)
# INTENTIONNEL : pas de préfixe '-' — si le fichier est absent, systemd refuse de démarrer
# (les credentials sont obligatoires, Program.cs fait aussi un fast-fail validation)
# Créer /etc/bet2invest-poster/env avec les variables suivantes :
#   Bet2Invest__Identifier=user@example.com
#   Bet2Invest__Password=your_password
#   Telegram__BotToken=123456:ABC-DEF...
#   Telegram__AuthorizedChatId=123456789
#   Poster__BankrollId=your_bankroll_id
# Protéger le fichier : chmod 600 /etc/bet2invest-poster/env, chown root:bet2invest
EnvironmentFile=/etc/bet2invest-poster/env

# Variables non-sensibles
Environment=DOTNET_ENVIRONMENT=Production

# Identification dans journald (journalctl -u bet2invest-poster)
SyslogIdentifier=bet2invest-poster

# Hardening sécurité
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/bet2invest-poster

[Install]
WantedBy=multi-user.target
